# templates/rrdcached-librenms.service.erb
[Unit]
Description=Data caching daemon for rrdtool
After=network.service

[Service]
Type=forking
ExecStart=/usr/bin/rrdcached -w 1800 -z 1800 -f 3600 -s librenms -U librenms -G librenms -j /opt/librenms/rrdcached/journal -l unix:/opt/librenms/rrdcached/rrdcached.sock -t 4 -F -b /opt/librenms/rrd/ -p /opt/librenms/rrdcached/rrdcached.pid
RRDC_USER=librenms

[Install]
WantedBy=default.target
